# deploy/systemd/control-api.service
[Unit]
Description=COITD Control API Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=cfst
Group=cfst
WorkingDirectory=/opt/cfst-collector
Environment=COLLECTOR_CONFIG_FILE=/opt/cfst-collector/config/collector.env
Environment=DDNS_CONFIG_FILE=/opt/cfst-collector/config/ddns.env
Environment=WEBUI_CONFIG_FILE=/opt/cfst-collector/config/webui.env
ExecStart=/usr/bin/python3 /opt/cfst-collector/scripts/control_api.py
Restart=always
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/opt/cfst-collector/data /opt/cfst-collector/logs /opt/cfst-collector/config

[Install]
WantedBy=multi-user.target
